# deploy/ansible/roles/snapshotter/files/goshimmer-db-dump.service
[Unit]
Description=Dump Goshimmer DB

[Service]
Type=oneshot
ExecStart=/opt/goshimmer_db_dump.sh

[Install]
WantedBy=timers.target
